# dist/socksy@.service
[Unit]
Description=Dirt simple SOCKS5 CONNECT proxy which binds a network interface
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/bin/socksy --bind-interface %i

[Install]
WantedBy=multi-user.target
